# bootstrap/system-config/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/autorandr --change --default default"
