# home/dot_config/systemd/user/rclone-sync.service
[Unit]
Description=Obsidian Bidirectional Sync (rclone bisync)
Documentation=man:rclone(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/dofs/obsidian-sync.env
ExecStart=/usr/bin/bash %h/dofs/scripts/obsidian/sync.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
