# home/.config/systemd/user/rclone-sync.timer
[Unit]
Description=Timer for Obsidian Bidirectional Sync

[Timer]
# Run 5 minutes after boot
OnBootSec=5min
# Run every 30 minutes thereafter
OnUnitActiveSec=30min
# Random delay to prevent hitting API limits if multiple machines start at once
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
